# dot_config/exact_sketchybar/sketchybar.ini
[program:sketchybar]
command = sketchybar
autostart = true 
autorestart = true
user = michael
redirect_stderr = true
environment = HOME="/Users/michael",XDG_CONFIG_HOME="/Users/michael/.config",XDG_DATA_HOME="/Users/michael/.local/share",PATH="/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin"
