# dot_config/exact_sketchybar/sketchybar.ini
[program:sketchybar]
command = /opt/homebrew/bin/sketchybar
autostart = true 
autorestart = true
user = michael
