# dot_config/exact_sketchybar/sketchybar.ini
[program:sketchybar]
command = /opt/homebrew/bin/sketchybar -c /Users/michael/.config/sketchybar/sketchybarrc
autostart = true 
autorestart = true
user = michael
redirect_stderr = true
