# dot_config/exact_sketchybar/sketchybar.ini
[program:sketchybar]
command=/opt/homebrew/bin/sketchybar -c /Users/michael/.config/sketchybar/sketchybarrc
autostart=true 
autorestart=true
user=michael
redirect_stderr=true
environment=HOME="/Users/michael",USER="michael",XDG_CONFIG_HOME="/Users/michael/.config",XDG_DATA_HOME="/Users/michael/.local/share",PATH="/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin"
